# sample/etc/system/systemd/onamaeddns.service.d/env.conf
[Service]
Environment="TARGET_HOST=<myhost>"
Environment="TARGET_DOMAIN=<example.com>"
